# packaging/debian/dd-trace-agent.upstart
description "dd-trace-agent"

# Defaults set by kernel
limit nofile 1024 4096

start on runlevel[2345]
stop on runlevel[016]

respawn limit 10 5

setuid dd-agent

exec /opt/datadog-agent/bin/trace-agent 2>&1
